# shell/.config/shell/profile.d/99-shell-path.sh
# shellcheck shell=sh

# Add our own utilities last so that they have precedence in the PATH.
PATH_try_prepend "$HOME/.local/bin"
